# crontab
#task
*/15 * * * * /home/yobit.sh
*/15 * * * * /home/micro.sh
